# _docker/shellvm/supervisord.conf
# Haelt den minimalen Desktop (Xvnc + openbox) dauerhaft am Leben — ersetzt
# das bisherige "sleep infinity" als CMD. Laeuft komplett als User "fallen"
# (Docker USER-Direktive, kein root hier). Logs -> stdout/stderr -> `docker
# logs glappa-shell-persistent`.
[supervisord]
nodaemon=true
loglevel=info
logfile=/var/log/glappa-desktop/supervisord.log
pidfile=/tmp/supervisord.pid

# -SecurityTypes None: bewusst OHNE eigenes VNC-Passwort. Der Schutz kommt
# aus der Netz-Isolation (glappa-shell-lan ist intern — nur der Gast selbst
# und displaygate koennen den Port ueberhaupt erreichen) UND aus displaygates
# eigenem Ticket-Login davor (siehe _docker/displaygate/). Ein zweites,
# separates VNC-Passwort haette hier nur Verwaltungsaufwand ohne echten
# zusaetzlichen Schutz bedeutet (klassische VNC-Auth ist zudem auf 8 Zeichen
# begrenzt — passt nicht zum eigentlichen Shell-Passwort).
# -localhost no: muss von AUSSERHALB dieses Containers erreichbar sein
# (displaygate laeuft in einem eigenen Container).
[program:xvnc]
command=/usr/bin/Xvnc :1 -geometry 1280x800 -depth 24 -SecurityTypes None -localhost no -AlwaysShared
autorestart=true
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Wartet bis Xvnc den X11-Unix-Socket angelegt hat (kann eine Idee langsamer
# hochkommen als openbox braucht, sonst "cannot open display :1").
[program:openbox]
command=sh -c 'until [ -e /tmp/.X11-unix/X1 ]; do sleep 0.2; done; DISPLAY=:1 exec openbox'
autorestart=true
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
