# _docker/shell-egress/supervisord.conf
# supervisord haelt die drei Egress-Dienste am Leben und startet sie bei
# Absturz neu (autorestart). Logs gehen nach stdout/stderr -> `docker logs
# glappa-shell-egress`. Bewusst kein Traffic-Logging in den Diensten selbst.
#
# Befehle bewusst OHNE absoluten Pfad -> supervisord sucht sie im PATH, das
# ist robust gegen /usr/bin vs /usr/sbin-Unterschiede zwischen Distributionen.
[supervisord]
nodaemon=true
user=root
loglevel=info

# tor laeuft als unprivilegierter tor-User (DataDirectory gehoert ihm, s.
# Dockerfile) — kein root-Warnhinweis, kein Ownership-Fehler.
[program:tor]
command=tor -f /etc/tor/torrc
user=tor
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# dnscrypt-proxy bindet :53 (privilegierter Port) -> laeuft als root. Der
# Container ist isoliert (kein docker.sock, keine Host-Mounts), das ist ok.
[program:dnscrypt-proxy]
command=dnscrypt-proxy -config /etc/dnscrypt-proxy/dnscrypt-proxy.toml
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:privoxy]
command=privoxy --no-daemon /etc/privoxy/config
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
